# sh.saqoo.sessylph.auto-adopt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Auto-adopt pipeline: checks for Claude Code updates daily at 09:00 JST -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>sh.saqoo.sessylph.auto-adopt</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>/Users/hiko/Documents/repos/Personal/Sessylph/scripts/auto-adopt.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/sessylph-auto-adopt.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sessylph-auto-adopt.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/hiko</string>
  </dict>
</dict>
</plist>
